# i3/.config/systemd/user/btops.service
[Unit]
Description=bspwm desktop manager

[Service]
ExecStart=/run/current-system/sw/bin/btops
Restart=always

[Install]
WantedBy=wm.target

