# home/.config/systemd/user/btops.service
[Unit]
Description=bspwm desktop manager

[Service]
ExecStart=/usr/bin/env btops
Restart=always

[Install]
WantedBy=wm.target

