# home/.config/systemd/user/dunst.service
[Unit]
Description=Dunst notification daemon

[Service]
ExecStart=/usr/bin/env dunst
Restart=always

[Install]
WantedBy=xsession.target

