# home/.config/systemd/user/polybar.target
[Unit]
Requires=polybar@eDP-1.service polybar@DP-1.service polybar@HDMI-2.service
WantedBy=xsession.target
