# home/.config/systemd/user/wm.target
[Unit]
Requires=xsession.target
BindsTo=xsession.target

